# smppsim/smppsim.path
[Unit]
Description=Monitor SMPPSim config file

[Path]
PathChanged=/vagrant/smppsim.props
