# kannel/smsbox.path
[Unit]
Description=Monitor Kannel config file

[Path]
PathChanged=/vagrant/smskannel.props
